# roles/jenkins/files/init_password.service
[Unit]
Description=Init Password for image
After=network.target 
Before=jenkins.service

[Service]
Type=oneshot
ExecStart=/bin/rm -rf /var/lib/jenkins/*
ExecStart=/bin/systemctl disable init-password

[Install]
WantedBy=multi-user.target
