# roles/jenkins/files/init_password.service
[Unit]
Description=Init Password for image
After=network.target 
Before=jenkins.service

[Service]
Type=oneshot
ExecStart=-/bin/apt install jenkins
ExecStart=-/bin/yum install jenkins
ExecStart=/bin/systemctl disable init-password

[Install]
WantedBy=multi-user.target
